# com.mactts.service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mactts.service</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Volumes/Store/Users/benjifc/.local/share/mactts/venv/bin/python3</string>
        <string>-m</string>
        <string>uvicorn</string>
        <string>main:app</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>8000</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Volumes/Store/Users/benjifc/.local/share/mactts</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Volumes/Store/Users/benjifc/.local/share/mactts/logs/service.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Volumes/Store/Users/benjifc/.local/share/mactts/logs/service.stderr.log</string>
</dict>
</plist>
